# service/create-guest-templates.service
[Unit]
Description=Import guest templates into XAPI
Requires=xapi.service xapi-init-complete.target
After=xapi.service xapi-init-complete.target network-init.service
ConditionPathExists=!/var/lib/misc/ran-create-guest-templates

[Service]
Type=oneshot
ExecStart=/usr/bin/create-guest-templates-wrapper
RemainAfterExit=yes
Restart=on-failure
RestartSec=60s

# Avoid tight restart loops if XAPI/import is persistently failing.
StartLimitIntervalSec=20min
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
